# root/etc/services.d/irssi/run
#!/usr/bin/with-contenv bash
screen -D -m -S \
        irssi s6-setuidgid abc /usr/bin/irssi \
        --home=/config/.irssi

sleep 1s
