# root/etc/services.d/flood/run
#!/usr/bin/with-contenv bash

cd /usr/flood || exit

exec \
	s6-setuidgid abc npm start
